# autoclone.service
[Unit]
Description=test

[Service]
ExecStart=/home/tgsachse/Programming/Python/AutoClone/acdaemon.py

[Install]
WantedBy=default.target
